# packages/flask-nginx/flask_nginx-0.8.15.tar.gz/flask_nginx-0.8.15/flask_nginx/templates/supervisord.ini
; -*- coding: utf-8 -*-
; vi:syntax=ini

[program:{{appname}}]
command={{venv}}/bin/celery -A {{celery|default(appname ~ '.celery')}} {% if gevent %} -P gevent {% endif %} \
    worker -l INFO -n {{appname}}@%%h --concurrency={{workers}}
process_name=%(program_name)s
user={{user}}
directory={{application_dir}}
numprocs=1
stopwaitsecs={{stopwait|default(30)}}               ; we may be processing something
redirect_stderr=true
stdout_logfile={{application_dir}}/instance/logs/celery-%(program_name)s.log
stopsignal=TERM
autorestart=true
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
environment=HOME="{{homedir}}",USER="{{user}}",PATH="{{venv}}/bin:{{miniconda}}:/usr/bin"
