# packages/flask-nginx/flask_nginx-0.8.15.tar.gz/flask_nginx-0.8.15/flask_nginx/templates/systemd.mount
[Unit]
Description=IRDS
[Mount]
What={{drive|default('//drive.irds.uwa.edu.au/sci-ms-001')}}
Where={{application_dir}}
Type=cifs
Options=_netdev,iocharset=utf8,rw,uid={{uid}},gid={{gid}},forceuid,forcegid,vers={{version|default('3.0')}},user={{user}},
    {%- if credentials is defined -%}
    credentials={{credentials}}
    {%- else -%}
    password={{password}}
    {%- endif %}
TimeoutSec=30

[Install]
WantedBy=multi-user.target
